# conf/supervisor.conf
[program:zeroecks]
command=zeroecks --config=/etc/zeroecks/config.py
autostart=true
